# etc/cron.hourly/reload-nginx
#!/bin/bash

# To reload certs
service nginx reload
